# runbot.service
[Unit]
Description=TG Bot service
After=network-online.target

[Service]
ExecStart=/home/python/run.sh

[Install]
WantedBy=multi-user.target
